# projects/sinister-os/source/iso-build/airootfs/etc/systemd/system/tailscaled.service.d/sinister.conf
# sinister.conf — Sinister OS drop-in override for tailscaled.service
# Author: RKOJ-ELENO :: 2026-05-24
#
# This is a systemd drop-in. It is loaded ON TOP of the upstream
# /usr/lib/systemd/system/tailscaled.service shipped by the `tailscale`
# pacman package. We do not own that base unit — we only patch behavior.
#
# What this override changes vs. upstream:
#   - Forces kernel-mode (TS_USERSPACE=false) for best LAN/exit-node throughput.
#   - Wires Wants= / After= network-online.target so the daemon doesn't race
#     the first DHCP lease and burn its initial backoff.
#   - Adds an ExecStartPost helper that auto-authenticates from
#     /etc/sinister/tailscale-authkey if the file is present and non-empty.
#     The helper is idempotent and self-deletes the key on a successful join.
#   - Bumps RestartSec= so a flaky tailnet doesn't tight-loop the daemon.
#
# Composes with:
#   - /usr/local/sbin/sinister-tailscale-firstboot.sh (the ExecStartPost target)
#   - /etc/sinister/tailscale-authkey (placeholder; operator drops real key)
#   - /etc/sinister/eve.toml (optional; supplies TS_EXTRA_ARGS)
#   - source/docker-stack/compose.mesh.yml (the FLEET-CONTAINERIZED sibling)

[Unit]
Wants=network-online.target
After=network-online.target

[Service]
Environment=TS_USERSPACE=false
Environment=TS_STATE_DIR=/var/lib/tailscale
RestartSec=10s
ExecStartPost=-/usr/local/sbin/sinister-tailscale-firstboot.sh
